# screenscripts/95-hdmi-plug.rules
ACTION=="change", SUBSYSTEM=="drm", ENV{DISPLAY}=":0", ENV{XAUTHORITY}="/home/ibb/.Xauthority", RUN+="/opt/smart-screens.sh"
